# scripts/bank-bot/launchd/cl.octopusfinance.bankbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>cl.octopusfinance.bankbot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/panchobriceno/Downloads/octopus-finance-source/scripts/bank-bot/run-daily.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/panchobriceno/Downloads/octopus-finance-source/scripts/bank-bot/logs/launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/panchobriceno/Downloads/octopus-finance-source/scripts/bank-bot/logs/launchd.err</string>
</dict>
</plist>
